# data/old_stuff/pi_receiver/system/ir_receiver.service.d/ir_receiver.conf
[Service]
Environment="ROOMNAME=HS1"
Environment="DOORNUMBER=1"
