# data/old_stuff/pi_emitter/system/sync_ip.timer
[Unit]
Description=Timer for data sync service

[Timer]
OnUnitActiveSec=60s
OnBootSec=20s

[Install]
WantedBy=timers.target
